# systemd/evoseal-watchdog.timer
[Unit]
Description=EVOSEAL Watchdog Timer (user)

[Timer]
OnBootSec=2m
OnUnitActiveSec=1m
Unit=evoseal-watchdog.service

[Install]
WantedBy=timers.target
