# scripts/evoseal.service
# EVOSEAL Systemd Service Template
# This template uses placeholders that are replaced during installation:
#   __EVOSEAL_ROOT__ - Full path to EVOSEAL installation directory
#   __EVOSEAL_USER__ - User account that will run the service
# Use install_evoseal_service.sh to install with proper substitutions

[Unit]
Description=EVOSEAL Evolution Service with Auto-Update
After=network.target
Requires=network.target

[Service]
Type=simple
User=__EVOSEAL_USER__
WorkingDirectory=__EVOSEAL_ROOT__

# Environment configuration
EnvironmentFile=-__EVOSEAL_ROOT__/.env
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=__EVOSEAL_ROOT__:__EVOSEAL_ROOT__/SEAL"
Environment="EVOSEAL_ROOT=__EVOSEAL_ROOT__"
Environment="EVOSEAL_VENV=__EVOSEAL_ROOT__/.venv"
Environment="EVOSEAL_LOGS=__EVOSEAL_ROOT__/logs"

# Configure logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=evoseal

# Core dumps (for debugging if needed)
LimitCORE=infinity

# Security options
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Resource limits
LimitNOFILE=65535
LimitNPROC=65535

# Main execution command
ExecStart=/bin/bash -c './scripts/evoseal-unified-runner.sh --mode=service'

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=300s
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
